# .config/autostart/radicale.desktop
[Desktop Entry]
Name=Radicale
Exec=bash -c "/home/doeringc/.config/radicale/start.sh"
Terminal=false
Type=Application
